# init/launch_agents/co.fluder.micropress.preview.plist
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>co.fluder.micropress.preview</string>
    <key>ProgramArguments</key>
    <array>
      <string>/usr/local/bin/micropress</string>
      <string>preview</string>
    </array>
    <key>KeepAlive</key>
    <true/>
  </dict>
</plist>
